# files/redis7479.ini
[program:redis_7479]
user=redis
command=/opt/redis/bin/redis-server /opt/redis/etc/redis7479.conf
autostart=true
autorestart=true
startsecs=1
startretries=9
redirect_stderr=false
stdout_logfile=/var/log/supervisor/redis7479.out
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_events_enabled=true
stderr_logfile=/var/log/supervisor/redis7479.err
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
stderr_events_enabled=true
